# router/etc/systemd/system/dump_ds.service
[Unit]
Description=Dump DS records

[Service]
Type=oneshot
ExecStart=/usr/bin/dump_ds_keys.sh

[Install]
WantedBy=multi-user.target
