# ns1/etc/systemd/system/dump_ds.timer
[Unit]
Description=Run DS dump every minute

[Timer]
OnBootSec=30
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
